# cakeagent.service
[Unit]
Description=CakeAgent — Personal AI Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cakeagent
Group=cakeagent
WorkingDirectory=/opt/cakeagent
ExecStart=/usr/bin/node dist/src/index.js
Restart=always
RestartSec=10

EnvironmentFile=/opt/cakeagent/.env

NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cakeagent /var/tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cakeagent

[Install]
WantedBy=multi-user.target
