# cakeagent.service
[Unit]
Description=CakeAgent — Personal AI Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cakeagent
Group=cakeagent
WorkingDirectory=/opt/cakeagent
ExecStart=/usr/bin/node dist/src/index.js
Restart=always
RestartSec=10

EnvironmentFile=/opt/cakeagent/.env

# Security: dedicated user, filesystem restrictions, kernel hardening
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# Agent workspace + apt package paths (sudo apt needs /usr, /var, /etc)
ReadWritePaths=/opt /usr /var /etc

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cakeagent

[Install]
WantedBy=multi-user.target
