# cakeagent.service
[Unit]
Description=CakeAgent — Personal AI Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cakeagent
Group=cakeagent
WorkingDirectory=/opt/cakeagent
ExecStart=/usr/bin/node dist/src/index.js
Restart=always
RestartSec=10
TimeoutStopSec=30

EnvironmentFile=/opt/cakeagent/.env

ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt /usr /var /etc

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cakeagent

[Install]
WantedBy=multi-user.target
