# cakeagent.service
[Unit]
Description=CakeAgent — Personal AI Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cakeagent
Group=cakeagent
WorkingDirectory=/opt/cakeagent
ExecStart=/usr/bin/node dist/src/index.js
Restart=always
RestartSec=10
TimeoutStopSec=30

EnvironmentFile=/opt/cakeagent/.env

# Filesystem: strict read-only baseline, scoped write access
ProtectSystem=strict
ReadWritePaths=/opt/cakeagent /usr /var /etc
InaccessiblePaths=/etc/shadow /etc/gshadow /etc/sudoers /etc/sudoers.d /etc/ssh

# Isolation
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
UMask=0077
RemoveIPC=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cakeagent

[Install]
WantedBy=multi-user.target
